# conf/supervisord/supervisord.conf
[supervisord]

[supervisorctl]

[program:countries]
command = /go/src/run

[program:nginx-app]
command = /usr/sbin/nginx
